# backend/gunicorn.service
[Unit]
Description=大颗粒项目后端
After=network.target

[Service]
User=vlight

Environment="MVCAM_COMMON_RUNENV=/opt/MVS/lib"
Environment="LD_LIBRARY_PATH=/home/vlight/dev/TensorRT-10.5.0.18/lib:/opt/MVS/lib/64:/opt/MVS/lib/32"

WorkingDirectory=/home/vlight/dev/big-particle/backend
ExecStart=/home/vlight/dev/big-particle/backend/.venv/bin/gunicorn
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
